# wmux@.service
[Unit]
Description=wmux - Web-based tmux Controller for %i
After=network-online.target tailscaled.service amux@%i.service
Wants=network-online.target amux@%i.service

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i
ExecStart=/usr/local/bin/wmux --tls --amux-url http://localhost:2023
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/tmp

# Environment
Environment="TERM=xterm-256color"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wmux

[Install]
WantedBy=multi-user.target
